# lab02/zad3/supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=mongod --dbpath /opt/app/data

[program:nodejs]
command=node server.js
directory=/opt/app
autostart=true
autorestart=true
stderr_logfile=/var/log/nodejs.err.log
stdout_logfile=/var/log/nodejs.out.log
